# supervisord.conf
[supervisord]
nodaemon=true

[program:h2db]
command=bash /h2-start.sh

[program:tomcat7]
command=catalina.sh run