# supervisord.conf
[supervisord]
nodaemon=true

[program:h2db]
command=bash /h2-start.sh
priority=1

[program:tomcat7]
command=catalina.sh run
priority=10